# templates/default/mongod.upstart.erb
# Generated by Chef - manual changes made will be lost!
description     "MongoDB server instance"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# If the process quits unexpectedly trigger a respawn
respawn

# Increase ulimit to recommended levels
limit nofile <%= node['mongodb']['open_file_limit'] %> <%= node['mongodb']['open_file_limit'] %>

# Start the process
exec start-stop-daemon --start --chuid <%= node['mongodb']['user'] %> --make-pidfile --pidfile /var/run/<%= @instance_name %>.pid --exec /usr/bin/mongod -- -f <%= @config_file %>

# Clean up after stopping the process
post-stop script
	rm /var/run/<%= @instance_name %>.pid
end script